# q_DBUS/kylin-container.service
[Unit]
Description=kylin-runing-server

[Service]
User=root
ExecStart=/usr/bin/kylin-run-service

[Install]
WantedBy=multi-user.target
