# shm-runing/kylin-container.service
[Unit]
Description=kylin-runing-server

[Service]
User=root
ExecStart=/usr/bin/shmread

[Install]
WantedBy=multi-user.target
