# msg-runing/kylin-container.service
[Unit]
Description=kylin-runing-server

[Service]
User=root
ExecStart=/usr/bin/msg-service

[Install]
WantedBy=multi-user.target
